# bar/packages/debian/postinst
#!/bin/sh

chmod 700 /etc/bar
chmod 600 /etc/bar/bar.cfg

if   test -d /lib/systemd/system; then
  install -d /lib/systemd/system
  install -m 644 /var/tmp/bar/lib/systemd/system/barserver.service /lib/systemd/system
fi
if   test -f /etc/SuSE-release -o -d /etc/SuSEconfig; then
  install -m 755 /var/tmp/bar/etc/init.d/barserver-SuSE /etc/init.d/barserver
elif test -f /etc/fedora-release; then
  install -m 755 /var/tmp/bar/etc/init.d/barserver-Fedora /etc/init.d/barserver"; \
elif test \( -f /etc/redhat-release -a -n "`grep 'Red Hat' /etc/redhat-release 2>/dev/null`" \); then
  install -m 755 /var/tmp/bar/etc/init.d/barserver-RedHat /etc/init.d/barserver"; \
elif test \( -f /etc/redhat-release -a -n "`grep 'CentOS' /etc/redhat-release 2>/dev/null`" \); then
  install -m 755 /var/tmp/bar/etc/init.d/barserver-CentOS /etc/init.d/barserver
elif test -f /etc/mandrake-release; then
  install -m 755 /var/tmp/bar/etc/init.d/barserver-Mandrake /etc/init.d/barserver
elif test -f /etc/lsb-release; then
  install -m 755 /var/tmp/bar/etc/init.d/barserver-debian /etc/init.d/barserver
elif test -f /etc/debian_release; then
  install -m 755 /var/tmp/bar/etc/init.d/barserver-debian /etc/init.d/barserver
else
  install -m 755 /var/tmp/bar/etc/init.d/barserver-debian /etc/init.d/barserver
fi

type chkconfig 1>/dev/null 2>/dev/null
if test $? -eq 0; then
  chkconfig --add barserver
  chkconfig barserver on
fi
type systemctl 1>/dev/null 2>/dev/null
if test $? -eq 0; then
  systemctl enable barserver
  systemctl start barserver
fi
type service 1>/dev/null 2>/dev/null
if test $? -eq 0; then
  service barserver start
fi

rm -rf /var/tmp/bar
