# bar/packages/debian/postinst
#!/bin/sh

chmod 700 /etc/bar
chmod 600 /etc/bar/bar.cfg

if test -d /lib/systemd; then
  if test ! -f /lib/systemd/system/barserver.service; then
    install -d /lib/systemd/system
    install -m 644 /var/tmp/bar/lib/systemd/system/barserver.service /lib/systemd/system
  fi
fi
install -m 755 /var/tmp/bar/etc/init.d/barserver-debian /etc/init.d/barserver

type chkconfig 1>/dev/null 2>/dev/null
if test $? -eq 0; then
  chkconfig --add barserver
  chkconfig barserver on
fi
type systemctl 1>/dev/null 2>/dev/null
if test $? -eq 0; then
  systemctl enable barserver
  systemctl start barserver
fi
type service 1>/dev/null 2>/dev/null
if test $? -eq 0; then
  service barserver start
fi
