# bar/osc/home:torsten20:BAR/BAR/debian.postinst
#!/bin/sh

chmod 700 /etc/bar
chmod 600 /etc/bar/bar.cfg

chkconfig --add barserver
chkconfig barserver on
